# services/claude-telegram-bridge/systemd/claude-telegram-bridge.service
[Unit]
Description=Claude Code Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude
WorkingDirectory=/opt/ella/claude-telegram-bridge
EnvironmentFile=/opt/ella/claude-telegram-bridge/.env
ExecStart=/opt/ella/claude-telegram-bridge/start-bridge.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
